# targets/arm64/openwrt/files/etc/init.d/dvpn-openwrt
#!/bin/sh /etc/rc.common
# Example script
# Copyright (C) 2007 OpenWrt.org

START=21

boot() {
	/usr/bin/dvpn-openwrt &
}
 
enable() {
        /usr/bin/dvpn-openwrt &
}

disable() {
	kill $(pidof /usr/bin/dvpn-openwrt)
}
